# scripts/.crontab
0 */4 * * * sh ~/development/push-my-post/scripts/cron_scan_sites.sh
0 */4 * * * sh ~/development/push-my-post/scripts/cron_pin_suggestions.sh